# first_install/system/desktop-updater.path
[Path]
PathChanged=/home/mopidy/.versionInfo

[Install]
WantedBy=multi-user.target
